# platforms/suse/apcupsd.in
#! /bin/sh
# Copyright (c) 1996 S.u.S.E. Gmbh Fuerth, Germany.  All rights reserved.
# Copyright (C) 1996-99 Andre M. Hedrick <andre@suse.com>
# Copyright (C) 1999-02 Riccardo Facchetti <riccardo@inetcosmos.org>
#
# Author: Andre Hedrick <hedrick@astro.dyer.vanderbilt.edu> 1996-99
#         Riccardo Facchetti <riccardo@inetcosmos.org> 1999-02
#
# /etc/init.d/apcupsd
#
### BEGIN INIT INFO
# Provides: apcupsd
# Required-Start: $network $syslog
# Required-Stop:  $network $syslog
# Default-Start:  2 3 5
# Default-Stop:
# Description:    Start the apcupsd daemon
### END INIT INFO

. /etc/rc.config

test "$START_APCUPSD" = yes || exit 0

return=$rc_done

APCPID=@PIDDIR@/apcupsd.pid
APCUPSD=@sbindir@/apcupsd

case "$1" in
    start)
        rm -f @PWRFAILDIR@/powerfail
        rm -f @nologdir@/nologin
        echo -n "Starting apcupsd power management"
        $APCUPSD || return=$rc_failed
        echo -e "$return"
    ;;
    stop)
        echo -n "Stopping apcupsd power management"
        if [ -f $APCPID ]; then
                kill $(cat $APCPID) || return=$rc_failed
                rm -f $APCPID
        else
                return=$rc_failed
        fi
        echo -e "$return"
    ;;
    restart)
        $0 stop
        $0 start
    ;;
    status)
        @sbindir@/apcaccess status
    ;;
    *)
    echo "Usage: $0 {start|stop|restart|status}"
    exit 1
esac

exit 0
